# src/novnc/conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -forever -shared -passwdfile /app/conf.d/passwd
autorestart=true
